# recipes-extended/wls/files/main-recovery.service
[Unit]
Description=Main Recovery

[Service]
Type=oneshot
ExecStart=/home/user/scripts/wlpro_main_recovery.sh
