# recipes-extended/wls/files/server-socket-recovery.service
[Unit]
Description=Server Socket Recovery

[Service]
Type=oneshot
ExecStart=/home/user/scripts/wlpro_server_socket_recovery.sh
